# airootfs/etc/xdg/autostart/apply-my-defaults.desktop
[Desktop Entry]
Type=Application
Name=Apply My Plasma Defaults
Exec=/usr/local/bin/apply-my-plasma-defaults
X-KDE-autostart-phase=1
NoDisplay=true
