# airootfs/etc/xdg/autostart/icon-fix.desktop
[Desktop Entry]
Type=Application
Name=Fix Installer Icon
Exec=cp /usr/share/applications/install-karch.desktop ~/Desktop/
X-KDE-autostart-phase=1